# public/service/laravel.ini
[program:swoole-http]
process_name=%(program_name)s_%(process_num)02d
command=/www/server/php/73/bin/php /www/wwwroot/wh.com/artisan swoole:http
autostart=true
autorestart=true
user=www
numprocs=1
redirect_stderr=true

[program:swoole-websocket]
process_name=%(program_name)s_%(process_num)02d
command=/www/server/php/73/bin/php /www/wwwroot/wh.com/artisan swoole:websocket
autostart=true
autorestart=true
user=www
numprocs=1
redirect_stderr=true

[program:swoole-dial]
process_name=%(program_name)s_%(process_num)02d
command=/www/server/php/73/bin/php /www/wwwroot/wh.com/artisan swoole:dial
autostart=true
autorestart=true
user=www
numprocs=1
redirect_stderr=true

[program:swoole-api]
process_name=%(program_name)s_%(process_num)02d
command=/www/server/php/73/bin/php /www/wwwroot/wh.com/artisan swoole:api
autostart=true
autorestart=true
user=www
numprocs=2
redirect_stderr=true

[program:esl-custom]
process_name=%(program_name)s_%(process_num)02d
command=/www/server/php/73/bin/php /www/wwwroot/wh.com/artisan esl:custom
autostart=true
autorestart=true
user=www
numprocs=1
redirect_stderr=true


